# scripts/linux/saki-panel-daemon.service
[Unit]
Description=Saki Panel Daemon
After=network.target saki-panel.service

[Service]
Type=simple
WorkingDirectory=/opt/saki-panel
Environment=NODE_ENV=production
Environment=DAEMON_HOST=127.0.0.1
Environment=DAEMON_PORT=24444
Environment=DAEMON_PANEL_URL=http://127.0.0.1:23333
ExecStart=/usr/bin/node apps/daemon/dist/index.js
Restart=always
RestartSec=5
User=webops
Group=webops

[Install]
WantedBy=multi-user.target
